# copy/opt/core/lib/svc/manifest/uptime-kuma.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest"
                name="export">
  <service name="application/uptime-kuma"
           type="service"
           version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network"
                grouping="require_all"
                restart_on="error"
                type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem"
                grouping="require_all"
                restart_on="error"
                type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <exec_method type="method"
                 name="start"
                 exec="node server/server.js &amp;"
                 timeout_seconds="60">
      <method_context working_directory="/opt/kuma">
        <method_credential user="kuma" group="kuma" />
        <method_environment>
          <envvar name="HOME" value="/opt/kuma" />
          <envvar name="DATA_DIR" value="/var/db/kuma" />
          <envvar name="UPTIME_KUMA_HOST" value="127.0.0.1" />
          <envvar name='PLAYWRIGHT_BROWSERS_PATH' value='/nonexistent'/>
          <envvar name="PATH" value="/bin:/usr/bin:/usr/local/bin:/opt/local/bin" />
        </method_environment>
      </method_context>
    </exec_method>
    <exec_method type="method"
                 name="stop"
                 exec=":kill"
                 timeout_seconds="60">
      <method_context/>
    </exec_method>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Uptime Kuma: A fancy self-hosted monitoring tool</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
